# deploy/launchd/com.prompt-preflight.proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Install: copy to ~/Library/LaunchAgents/ then:
     launchctl load ~/Library/LaunchAgents/com.prompt-preflight.proxy.plist -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.prompt-preflight.proxy</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/enhance</string>
    <string>--serve-only</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PROMPT_ENHANCER_BACKEND</key>
    <string>auto</string>
  </dict>
  <key>StandardErrorPath</key>
  <string>/tmp/prompt-preflight-proxy.log</string>
</dict>
</plist>
